# turnkey_web/tkw.ini
[program:tkw]
# TCP socket used by Nginx backend upstream
#socket=tcp://localhost:8001
chown=user:user

# Directory where your site's project files are located
directory=/home/user/django-taiwan-einvoice/turnkey_web

# Each process needs to have a separate socket file, so we use process_num
# Make sure to update "mysite.asgi" to match your project name
environment=TURNKEY_WEB_LOG_LEVEL="DEBUG",ALLOWED_HOSTS="tkw.example.com,ttkw.example.com",DEFAULT_DB_HOST="deti-db.ho600.com"

command=/home/user/django-taiwan-einvoice.py3env/bin/daphne -b 0.0.0.0 -e ssl:443:privateKey=./django-taiwan-einvoice.key:certKey=./django-taiwan-einvoice.crt turnkey_web.asgi:application

# Number of processes to startup, roughly the number of CPUs you have
numprocs=1

# Give each process a unique name so they can be told apart
process_name=asgi%(process_num)d

# Automatically start and recover processes
autostart=true
autorestart=true

# Choose where you want your log to go
stdout_logfile=/home/user/django-taiwan-einvoice/turnkey_web.asgi.log
redirect_stderr=true
